# roles/app/templates/app.service.j2
[Unit]
Description=App xpaste_practicum

[Service]
Type=simple
ExecStart=/opt/xpaste_practicum/bin/start.sh

[Install]
WantedBy=default.target

